# cloudflared/cf
#!/sbin/openrc-run

# supervisor=supervisor-daemon
name=cf

command=/bin/cloudflared
command_args="tunnel run"

pidfile="/var/run/cf.pid"
respawn_delay=5
respawn_max=0

set -o allexport
if [ -f /etc/environment ]; then source /etc/environment; fi
set +o allexport
